# ai-server/llama-swap.service
[Unit]
Description=LLaMA Swap Service
After=network.target

[Service]
Type=simple
ExecStart=/home/crussell/Code/llama-swap -listen 0.0.0.0:8000
WorkingDirectory=/home/crussell/Code
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
