# servers/hub/openclaw/openclaw.container.in
[Unit]
Description=OpenClaw gateway
After=network-online.target podman.socket
Wants=network-online.target podman.socket
Requires=podman.socket

[Container]
Image=localhost/openclaw-hub:latest
ContainerName=openclaw
UserNS=keep-id
User=%U:%G
SecurityLabelDisable=true
DropCapability=all
NoNewPrivileges=true
RunInit=true
PidsLimit=2048
Memory=8g
PublishPort=127.0.0.1:30110:18789
PublishPort=127.0.0.1:30111:18790
Volume=/srv/openclaw/config:/home/node/.openclaw:Z
Volume=/srv/openclaw/workspace:/srv/openclaw/workspace:Z
Volume=/srv/openclaw/sandboxes:/srv/openclaw/sandboxes:Z
Volume={{PODMAN_SOCKET}}:/var/run/docker.sock
Environment=HOME=/home/node
Environment=TERM=xterm-256color
Environment=DOCKER_HOST=unix:///var/run/docker.sock
EnvironmentFile=/srv/openclaw/config/.env
Pull=never
Exec=node /app/openclaw.mjs gateway --bind lan --port 18789

[Service]
TimeoutStartSec=900
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
