# servers/hub/dev-stacks/buildspace/systemd/units/buildspace.target
[Unit]
Description=All Buildspace services
Wants=buildspace-stack.target
After=buildspace-stack.target
