# servers/hub/dev-stacks/buildspace/systemd/units/buildspace-marketplace.service
[Unit]
Description=Buildspace marketplace app
After=buildspace-infra.target buildspace-deps.service
Wants=buildspace-infra.target buildspace-deps.service
ConditionPathExists=%h/Code/bs/buildspace/package.json
ConditionPathExists=%h/Code/bs/buildspace/.env
PartOf=buildspace-stack.target

[Service]
Type=simple
WorkingDirectory=%h/Code/bs/buildspace
EnvironmentFile=%h/.config/systemd/user/buildspace-overrides.env
ExecStart=%h/.bun/bin/bun --env-file=.env --filter @buildspace/marketplace dev
Restart=on-failure
RestartSec=2

[Install]
WantedBy=buildspace-stack.target
