# servers/hub/dev-stacks/buildspace/systemd/units/buildspace-bootstrap.service
[Unit]
Description=Bootstrap Buildspace local database
After=buildspace-infra.target
Wants=buildspace-infra.target
ConditionPathExists=%h/Code/bs/buildspace/package.json
ConditionPathExists=%h/Code/bs/buildspace/.env

[Service]
Type=oneshot
WorkingDirectory=%h/Code/bs/buildspace
EnvironmentFile=%h/.config/systemd/user/buildspace-overrides.env
ExecStart=%h/.bun/bin/bun install --frozen-lockfile
ExecStart=%h/.bun/bin/bun --env-file=.env db:migrate
ExecStart=%h/.bun/bin/bun --env-file=.env db:seed
TimeoutStartSec=30min
