# servers/hub/dev-stacks/buildspace/systemd/units/buildspace-infra.target
[Unit]
Description=Buildspace shared infra
# Start compose postgres separately: just infra-up
