# datenight/datenight.container
[Unit]
Description=Date Night Restaurant Picker
After=network-online.target
Wants=network-online.target

[Container]
Image=localhost/datenight:latest
ContainerName=datenight
PublishPort=30790:7890
Volume=/srv/datenight/data:/data:Z
Environment=DATA_FILE=/data/restaurants.json

[Service]
Restart=unless-stopped
RestartSec=5

[Install]
WantedBy=default.target
