# servers/hub/backup/var-mnt-tank-backups.automount
[Unit]
Description=Automount for NAS backups

[Automount]
Where=/var/mnt/tank/backups
TimeoutIdleSec=300

[Install]
WantedBy=multi-user.target
