# brunch/config/modules/dev-infra/dev-postgres.container
[Unit]
Description=Shared dev Postgres
After=dev-network.service
Requires=dev-network.service dev-postgres-data-volume.service

[Container]
ContainerName=dev-postgres
Image=docker.io/library/postgres:16
Pull=missing
Network=dev.network
NetworkAlias=postgres
PublishPort=127.0.0.1:5432:5432
Environment=POSTGRES_USER=postgres
Environment=POSTGRES_PASSWORD=postgres
Environment=POSTGRES_DB=postgres
HealthCmd=pg_isready -U postgres -d postgres
HealthInterval=15s
HealthRetries=20
HealthStartupInterval=5s
HealthStartupRetries=30
HealthTimeout=5s
Notify=healthy
Volume=dev_postgres_data:/var/lib/postgresql/data

[Service]
Restart=always
