# servers/hub/quadlets/containers/beszel-agent.container
[Unit]
Description=Beszel Agent for crussell-srv
After=network-online.target
Wants=network-online.target

[Container]
Image=docker.io/henrygd/beszel-agent:latest
Network=host
Volume=/srv/beszel/agent:/var/lib/beszel-agent
Environment=LISTEN=45876
Environment=HUB_URL=http://127.0.0.1:8090
EnvironmentFile=/srv/beszel/secrets.env

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
