# servers/hub/quadlets/containers/mastra.container
[Unit]
Description=Mastra Studio and API Server
After=network-online.target
Wants=network-online.target

[Container]
Image=localhost/mastra:latest
PublishPort=30411:4111
Volume=/srv/mastra/data:/data
Environment=PORT=4111
Environment=OPENAI_BASE_URL=http://192.168.20.41:8000/v1
Environment=OPENAI_API_KEY=local-llama
Environment=MASTRA_MODEL=qwen3.5-9b
Environment=MASTRA_DB_URL=file:/data/mastra.db
Environment=MASTRA_TELEMETRY_DISABLED=1
HealthCmd=/usr/bin/wget --quiet --tries=1 --spider http://localhost:4111/health
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=40s

[Service]
Restart=unless-stopped
RestartSec=5

[Install]
WantedBy=default.target
