# servers/hub/dev-stacks/gloo/systemd/units/gloo-hummingbird.target
[Unit]
Description=Gloo Hummingbird dev stack
Wants=gloo-hb-api.service gloo-hb-web.service
