# servers/hub/dev-stacks/gloo/systemd/units/gloo-storyhub.target
[Unit]
Description=Gloo Storyhub dev stack
Wants=gloo-storyhub.service gloo-storyhub-worker.service
