# bee/openclaw-node/openclaw-node.service
[Unit]
Description=OpenClaw Node Daemon (Podman)
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
NotifyAccess=all
ExecStart=/usr/bin/podman run --name openclaw-node \
  --rm \
  --conmon-pidfile=%t/openclaw-node.pid \
  --cgroup-manager=systemd \
  --sdnotify=conmon \
  -e OPENCLAW_MODE=node \
  -e OPENCLAW_GATEWAY_URL=ws://192.168.20.32:30087 \
  -e NODE_ENV=production \
  -v /var/home/crussell/.openclaw-node/data:/root/.openclaw \
  -v /var/home/crussell/Code:/root/workspace:ro \
  --network=host \
  ghcr.io/openclaw/openclaw:latest

ExecStop=/usr/bin/podman stop -t 10 openclaw-node

Restart=always
RestartSec=10s
TimeoutStartSec=180s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/var/home/crussell/.openclaw-node /var/home/crussell/.openclaw

[Install]
WantedBy=multi-user.target
