# crussell-srv/quadlets/containers/immich-postgres.container
[Unit]
Description=Immich Postgres Database
After=network-online.target
Wants=network-online.target

[Container]
Image=docker.io/pgvector/pgvector:pg16
ContainerName=immich-postgres
Environment=POSTGRES_DB=immich
Environment=POSTGRES_USER=immich
Environment=POSTGRES_PASSWORD=immich_secret_password_change_me
Volume=immich-postgres-data:/var/lib/postgresql/data
Network=immich.network

[Service]
Restart=unless-stopped
ExecStartPost=/usr/bin/pg_isready -U immich -d immich

[Install]
WantedBy=default.target
